# automation-report/launchd/com.thiennp.automation-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.thiennp.automation-report</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npm</string>
    <string>run</string>
    <string>start</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/thien.nguyen/thiennp.github.io/automation-report</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>AUTOMATION_REPORT_HOST</key>
    <string>127.0.0.1</string>
    <key>AUTOMATION_REPORT_PORT</key>
    <string>3120</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/thien.nguyen/thiennp.github.io/automation-report/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/thien.nguyen/thiennp.github.io/automation-report/logs/launchd.err.log</string>
</dict>
</plist>
