# .codex-tmp/agent-report/launchd/com.thiennp.agent-report.autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.thiennp.agent-report.autostart</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>AGENT_REPORT_HOST=127.0.0.1 /Users/thiennguyen/thiennp.github.io/agent-report/scripts/ensure-agent-report-server.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/thiennguyen/thiennp.github.io/agent-report/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/thiennguyen/thiennp.github.io/agent-report/logs/launchd.err.log</string>
</dict>
</plist>
